# rootfs/etc/services.d/dnsmasq/run
#!/command/with-contenv sh
set -e -u -o pipefail

while ! $(ip link | grep -q tun0); do sleep 1; done

echo "Starting dnsmasq..."
dnsmasq --no-daemon --log-queries

sleep infinity &

wait